# packaging/prerm
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
  systemctl stop ambipi.service || true
  if [ "$1" = "remove" ]; then
    systemctl disable ambipi.service || true
  fi
fi

exit 0
